# overlay-debian-bullseye/etc/udev/rules.d/70-backlight.rules
KERNEL=="openframe-bl", RUN+="/bin/chmod 666 /sys/class/backlight/openframe-bl/bl_power /sys/class/backlight/openframe-bl/brightness"
