# overlay-ubuntu-bionic/lib/systemd/system/of-update.service
[Unit]
Description=OpenFrame Update

[Service]
ExecStart=/usr/local/sbin/of-update systemd
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
