# overlay-trusty/etc/network/if-pre-up.d/openframe-mac
#!/bin/bash

# openframe-mac v1.01 (7th November 2014)
#  Manages MAC addresses on OpenFrame devices.
#  We copy the MAC from wlan0 to eth0, but with the local bit set.
#  If no wlan0 interface is available, we generate a MAC for the eth0 interface.

if [ "$IFACE" = "wlan0" ]; then

	WIFIMAC="`ip link show $IFACE | grep "link/ether" | cut -d" " -f6`"
	[ "$WIFIMAC" != "" ] && echo $WIFIMAC > /tmp/mac_$IFACE

fi

if [ "$IFACE" = "eth0" ]; then

	COUNT=0
	while [ $COUNT -lt 10 ]; do
		if [ -f /tmp/mac_wlan0 ]; then
			WIFIMAC=`cat /tmp/mac_wlan0`
			SOMEOCTETS=`echo $WIFIMAC | cut -b 4-`
			ETHMAC="02:$SOMEOCTETS"
			ip link set $IFACE address $ETHMAC
			rm /etc/udev/rules.d/70-persistent-net.rules 2>/dev/null
			echo "# Automatically generated by /etc/network/if-pre-up.d/openframe-mac. Edits will be lost." > /etc/udev/rules.d/70-persistent-net.rules
			echo "SUBSYSTEM==\"net\", ACTION==\"add\", DRIVERS==\"?*\", ATTR{address}==\"$ETHMAC\", ATTR{dev_id}==\"0x0\", ATTR{type}==\"1\", KERNEL==\"eth*\", NAME=\"eth0\"" >> /etc/udev/rules.d/70-persistent-net.rules
			echo "SUBSYSTEM==\"net\", ACTION==\"add\", DRIVERS==\"?*\", ATTR{address}==\"$WIFIMAC\", ATTR{dev_id}==\"0x0\", ATTR{type}==\"1\", KERNEL==\"wlan*\", NAME=\"wlan0\"" >> /etc/udev/rules.d/70-persistent-net.rules
			echo "Adaptation of wlan0 MAC for eth0 successful after $COUNT seconds." > /tmp/mac_status_$IFACE
			exit 0
		else
			sleep 1
			let COUNT+=1
		fi
	done

	NOWIFIMAC="`grep "eth0" /etc/udev/rules.d/70-persistent-net.rules | grep -o -E "([0-9a-fA-F]{2}(\:|\-)){5}[0-9a-fA-F]{2}"`" || true
	if [ "$NOWIFIMAC" = "" ]; then
		NOWIFIMAC="02:`cat /proc/interrupts | md5sum | sed -r "s/(.{2})/\1:/g; s/^(.{14}).*/\1/"`"
		echo "No wlan0 found after $COUNT seconds. Using newly generated $NOWIFIMAC address." > /tmp/mac_status_$IFACE
		rm /etc/udev/rules.d/70-persistent-net.rules 2>/dev/null
		echo "# Automatically generated by /etc/network/if-pre-up.d/openframe-mac. Edits will be lost." > /etc/udev/rules.d/70-persistent-net.rules
		echo "SUBSYSTEM==\"net\", ACTION==\"add\", DRIVERS==\"?*\", ATTR{address}==\"$NOWIFIMAC\", ATTR{dev_id}==\"0x0\", ATTR{type}==\"1\", KERNEL==\"eth*\", NAME=\"eth0\"" >> /etc/udev/rules.d/70-persistent-net.rules
	else
		echo "No wlan0 found after $COUNT seconds. Using $NOWIFIMAC address from rules." > /tmp/mac_status_$IFACE
	fi

	ip link set $IFACE address $NOWIFIMAC

fi

exit 0
