# overlay-bionic-rescue/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -a OPENFRAMEUSER --noclear %I $TERM
