# overlay-bionic/etc/network/if-up.d/ntpkick
#!/bin/bash
ntpdate 0.pool.ntp.org &>/dev/null
ntpdate 1.pool.ntp.org &>/dev/null
