# overlay-trusty/etc/network/if-up.d/ntpkick
#!/bin/bash
ntpdate ntp2b.mcc.ac.uk &>/dev/null
ntpdate ntp2d.mcc.ac.uk &>/dev/null
